# files/pacman-dl.timer
[Unit]
Description=triggers pacman-dl every 4 to 6 hours.

[Timer]
OnCalendar=00/5:0
RandomizedDelaySec=60min

[Install]
WantedBy=timers.target
